# packages/mcp-server/com.lokascript.mcp-server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lokascript.mcp-server</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update path to your node installation -->
        <string>/opt/homebrew/opt/node@20/bin/node</string>
        <!-- Update path to your lokascript installation -->
        <string>$LOKASCRIPT_PATH/packages/mcp-server/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>$LOKASCRIPT_PATH/packages/mcp-server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>$LOKASCRIPT_PATH/packages/mcp-server/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>$LOKASCRIPT_PATH/packages/mcp-server/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
